# upstart/heysaopaulo.conf
description "heysaopaulo"

start on (filesystem)
stop on runlevel [016]

respawn
console log
setuid hey
setgid hey
chdir /home/hey/heysaopaulo

exec /home/hey/heysaopaulo-env/bin/gunicorn_django
